# .config/systemd/user/foot.service
[Unit]
Description=Foot

[Service]
ExecStart=/usr/bin/foot --server
Type=simple

[Install]
WantedBy=sway-session.target
