# .config/systemd/user/sway-session.target
[Unit]
BindsTo=graphical-session.target
